# bot/systemd/claude-worker@.service
[Unit]
Description=Claude Code GLM Worker %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claudebot
WorkingDirectory=/srv/claude-bot
EnvironmentFile=/etc/claude-bot/claude-bot.env
ExecStart=/usr/bin/python3 /srv/claude-bot/worker.py --worker-id %i
Restart=always
RestartSec=10

CPUQuota=250%
MemoryHigh=4G
MemoryMax=6G
TasksMax=300

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/srv/claude-bot /tmp /home/claudebot/.claude

[Install]
WantedBy=multi-user.target
