# bot/systemd/claude-patrol.timer
[Unit]
Description=Claude Bot Patrol Timer
Requires=claude-patrol.service

[Timer]
OnBootSec=30min
OnUnitActiveSec=6h
RandomizedDelaySec=15min
AccuracySec=1min

[Install]
WantedBy=timers.target
